# gvrepsrv.conf
[program:gvrepsrv]
command = /abspath/gvrepsrv --port=8000
directory = /abspath
stdout_logfile = gvrepsrv.log
redirect_stderr = true ; 2>&1
autorestart = true
